# com.airpodsmaxbattery.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Edit the path below to wherever the repo lives on your machine. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.airpodsmaxbattery</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/airpods-max-battery/AirPodsMaxBattery.app/Contents/MacOS/AirPodsMaxBattery</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
